# EIB/os-files/etc/systemd/system/elemental-register-systemd.service
[Unit]
Description=Elemental Register Install via Systemd
Wants=network-online.target
After=network-online.target
# Wait for hostname to be set (combustion script sets it)
# Check that hostname is not the default 'slemicro' or empty
ConditionPathExists=!/etc/rancher/elemental/agent/elemental_connection.json

[Install]
WantedBy=network-online.target

[Service]
EnvironmentFile=-/etc/sysconfig/proxy
Type=oneshot
# Use wrapper script to wait for hostname before registering
# This ensures ${Runtime/Hostname} template variable is resolved correctly
ExecStart=/usr/local/bin/wait-for-hostname-and-register.sh
ExecStartPost=/usr/bin/cp /var/lib/elemental/agent/elemental_connection.json /etc/rancher/elemental/agent
Restart=on-failure
RestartSec=10

